# bpn/crontabs/unconfirmed.sh
#!/bin/bash
date

php /var/www/bitping/monitor/unconfirmed-monitor.php &
sleep 10

php /var/www/bitping/monitor/unconfirmed-monitor.php &
sleep 10  

php /var/www/bitping/monitor/unconfirmed-monitor.php &
sleep 10  

php /var/www/bitping/monitor/unconfirmed-monitor.php &
sleep 10  

php /var/www/bitping/monitor/unconfirmed-monitor.php &
sleep 10  

php /var/www/bitping/monitor/unconfirmed-monitor.php &
